# vaxp-session.target
[Unit]
Description=VAXP User Session
BindsTo=graphical-session.target
Wants=graphical-session-pre.target
After=graphical-session-pre.target
